# etc/crontab
*/15 * * * * venv/bin/python ./backend/secret_cleaner.py